# pi-radio.service
[Unit]
Description=Pi Radio - Internet Radio Controller
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%PROJECT_DIR%
ExecStartPre=/bin/sleep 5
ExecStartPre=/usr/bin/pulseaudio --check || /usr/bin/pulseaudio --start
ExecStart=%PROJECT_DIR%/.venv/bin/python %PROJECT_DIR%/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
